# supervisor/hurl.conf
[program:wildcard.hurl]
command=bundle exec shotgun config.ru
directory=/home/www/wildcard.hurl/hurl
user=www
autostart=true
autorestart=true
redirect_stderr=True

