# supervisor/pad.wildcard.conf
[program:pad.wildcard]
command=/home/www/pad.wildcard/etherpad-lite/bin/run.sh
directory=/home/www/pad.wildcard/etherpad-lite
user=www
autostart=true
autorestart=true
redirect_stderr=True

